# rr-addons-main/eudev/src/99-qed.rules
SUBSYSTEM=="pci", DRIVER=="qed*", ACTION=="change", ENV{QED_DEBUGFS_BDF_DBG}=="?*", RUN+="/lib/udev/qed_udev_dbg.sh"
